# config/supervisor/elasticsearchHead.ini
[program:elasticsearchHead]
directory=/elasticsearchHead
command=/usr/bin/npm start
stdout_logfile=/logs/elasticsearchHead.log
redirect_stderr=true
autostart=true