# config/supervisor/elasticsearch.ini
[program:elasticsearch]
command=/elasticsearch/bin/elasticsearch
stdout_logfile=/logs/elasticsearch/supervisor.log
redirect_stderr=true
autostart=true
;directory=
user=elasticsearch
