# config/supervisor/kibana.ini
[program:kibana]
command = /kibana/bin/kibana
stdout_logfile=/logs/kibana/supervisor.log
stopasgroup = true
umask=022
priority=30
autostart=true
autorestart=true
startsecs=8
startretries=3
stopsignal=KILL
stopwaitsecs=5
user=kibana
redirect_stderr=true