# config/supervisor/logstash.ini
[program:logstash]
command=/logstash/bin/logstash agent -f /logstash/etc/
redirect_stderr=true
autostart=true
stdout_logfile=/logs/logstash/supervisor.log